# misc/systemd/cf-consumer.service.in
[Unit]
Description=CFEngine Enterprise Report Consumer
PartOf=cfengine3.service
After=syslog.target
ConditionPathExists=@workdir@/bin/cf-consumer

After=cf-postgres.service
BindsTo=cf-postgres.service
After=cf-redis-server.service
BindsTo=cf-redis-server.service

[Service]
Type=simple
ExecStart=@workdir@/bin/cf-consumer --no-fork
Restart=always
RestartSec=10
PIDFile=@workdir@/cf-consumer.pid

[Install]
WantedBy=multi-user.target
WantedBy=cfengine3.service
WantedBy=cf-redis-server.service
