# misc/systemd/cf-execd.service.in
[Unit]
Description=CFEngine Enterprise Execution Scheduler
After=syslog.target
ConditionPathExists=@workdir@/bin/cf-execd
ConditionPathExists=@workdir@/inputs/promises.cf
PartOf=cfengine3.service

[Service]
Type=simple
ExecStart=@workdir@/bin/cf-execd --no-fork
Restart=always
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target
WantedBy=cfengine3.service
