# misc/systemd/cf-hub.service.in
[Unit]
Description=CFEngine Enterprise Hub Report Collector
PartOf=cfengine3.service
ConditionPathExists=@workdir@/bin/cf-hub
ConditionPathExists=@workdir@/inputs/promises.cf
After=syslog.target
After=network.target

Wants=cf-postgres.service
After=cf-postgres.service

[Service]
Type=simple
ExecStart=@workdir@/bin/cf-hub  --no-fork
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
WantedBy=cfengine3.service
