# misc/systemd/cf-monitord.service.in
[Unit]
Description=CFEngine Monitor Daemon
After=syslog.target
ConditionPathExists=@bindir@/cf-monitord
ConditionPathExists=@workdir@/inputs/promises.cf
PartOf=cfengine3.service

[Service]
Type=simple
ExecStart=@bindir@/cf-monitord --no-fork
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
WantedBy=cfengine3.service
