# misc/systemd/cf-postgres.service.in
[Unit]
Description=CFEngine Enterprise PostgreSQL Database
After=syslog.target
ConditionPathExists=@workdir@/bin/pg_ctl
PartOf=cfengine3.service

[Service]
Type=forking
WorkingDirectory=/tmp
User=cfpostgres
ExecStart=@workdir@/bin/pg_ctl -w -D @statedir@/pg/data -l /var/log/postgresql.log start
ExecStop=@workdir@/bin/pg_ctl -w -D @statedir@/pg/data -l /var/log/postgresql.log stop -m fast

ExecReload=@workdir@/bin/pg_ctl -w -D @statedir@/pg/data -l /var/log/postgresql.log reload -m fast
PIDFile=@statedir@/pg/data/postmaster.pid
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
WantedBy=cfengine3.service
