# misc/systemd/cf-redis-server.service.in
[Unit]
Description=CFEngine Enterprise Redis Report Collection Cache
After=syslog.target
PartOf=cfengine3.service
ConditionPathExists=@workdir@/bin/redis-server
Wants=cf-consumer.service
Before=cf-hub.service
Wants=cf-hub.service
Before=cf-consumer.service

[Service]
Type=forking
ExecStart=@workdir@/bin/redis-server @workdir@/config/redis.conf
PIDFile=@workdir@/redis-server.pid
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
WantedBy=cfengine3.service
WantedBy=cf-hub.service
