# misc/systemd/cf-serverd.service.in
[Unit]
Description=CFEngine Enterprise file server
After=syslog.target
After=network.target
ConditionPathExists=@bindir@/cf-serverd
ConditionPathExists=@workdir@/policy_server.dat
ConditionPathExists=@workdir@/inputs/promises.cf
PartOf=cfengine3.service

[Service]
Type=simple
ExecStart=@bindir@/cf-serverd --no-fork
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
WantedBy=cfengine3.service
